# .config/systemd/user/npm-ai-update.timer
[Unit]
Description=Daily AI npm package update check

[Timer]
OnCalendar=daily
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
